# plugins/workflow-engine/systemd/hermes-workflow-dispatcher.service
# NOTE: A future upstream PR will add workflow.dispatch_in_gateway (mirroring
# kanban.dispatch_in_gateway) so the scheduler runs inside the gateway process.
# When that ships, disable this unit and rely on the gateway instead:
#
#     systemctl --user disable --now hermes-workflow-dispatcher.service
#     # ensure a gateway is running (e.g. via hermes-gateway.service)
#
# This unit is the Option 2 standalone daemon; keep it for environments where
# running a long-lived gateway process is not possible.
# Running this unit AND a future gateway-embedded scheduler simultaneously
# is NOT supported — they would race against the same switchui-workflows.db
# (the cross-process file lock in engine/db/migrate.py serialises schema
# migrations only; double-dispatch is a logic error, not a crash).

[Unit]
Description=Hermes Workflow Engine dispatcher (standalone daemon — see workflow.dispatch_in_gateway for future gateway-embedded mode)
Documentation=https://hermes-agent.nousresearch.com/docs/user-guide/features/workflow-engine
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env hermes workflow daemon --interval 60 --pidfile %t/hermes-workflow-dispatcher.pid
Restart=on-failure
RestartSec=5
# Log to the journal via stdout/stderr
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
